# deploy/pi5/cortex-watchdog.service
[Unit]
Description=Cortex daemon health watchdog probe
Documentation=file:///etc/rudy/docs/runbooks/operator-console.md
After=network-online.target cortex.service
Wants=network-online.target cortex.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cortex-watchdog.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
