# deploy/pi5/cortex-watchdog.timer
[Unit]
Description=Probe cortex health every 15 seconds
Documentation=file:///etc/rudy/docs/runbooks/operator-console.md

[Timer]
OnBootSec=60s
OnUnitActiveSec=15s
AccuracySec=2s
Unit=cortex-watchdog.service
Persistent=true

[Install]
WantedBy=timers.target
