# deploy/pi5/cortex-update.timer
[Unit]
Description=Poll GitHub for new Cortex Pi releases (every minute)
Documentation=file:///etc/rudy/docs/runbooks/operator-console.md

[Timer]
# First poll 30s after boot, then every 60s.
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=cortex-update.service
Persistent=true

[Install]
WantedBy=timers.target
