# deploy/pi5/rudy-can.service
[Unit]
Description=Rudy SocketCAN bringup (Waveshare 2-CH HAT)
After=network-pre.target
Before=network-online.target rudy-driver.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/rudy-can-setup.sh 1000000

[Install]
WantedBy=multi-user.target
